# linux/platform/common/fs/buildroot/fs/skeleton/etc/profile.d/S02SDL.sh
#!/bin/sh
export SDL_MOUSEDRV=TSLIB
export SDL_MOUSEDEV=$TSLIB_TSDEVICE

